# base/drydock-desks.service
# Drydock on-boot desk resume — systemd system unit (Linux).
# Brings every workspace in state=suspended back to running after a reboot.
# Pairs with drydock-wsd.service; depends on wsd being ready.
#
# Install alongside drydock-wsd.service:
#   cp base/drydock-desks.service /etc/systemd/system/
#   cp scripts/drydock-resume-desks /usr/local/bin/
#   chmod +x /usr/local/bin/drydock-resume-desks
#   systemctl daemon-reload
#   systemctl enable drydock-desks.service
#
# Oneshot + RemainAfterExit so `systemctl status drydock-desks` shows
# the last run's outcome, and a manual `systemctl start drydock-desks`
# re-runs the resume pass.
[Unit]
Description=Drydock — resume suspended desks at boot
Documentation=https://github.com/stevefan/drydock
After=drydock-wsd.service
Requires=drydock-wsd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/drydock-resume-desks
# On graceful shutdown: mark running desks suspended so they resume next boot.
# drydock-desks stops before drydock-wsd (reverse of start order), so wsd is
# still serving when ws stop runs.
ExecStop=/usr/local/bin/drydock-stop-desks
TimeoutStopSec=180
StandardOutput=append:/root/.drydock/logs/desks-resume.log
StandardError=append:/root/.drydock/logs/desks-resume.log
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/root/.local/bin
Environment=HOME=/root

[Install]
WantedBy=multi-user.target
